# gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket
[Socket]
ListenStream=/mnt/ebs/gunicorn/gunicorn.sock
[Install]
WantedBy=sockets.target